# os_maintain/system-patch/etc/udev/rules.d/vga_switch.rules
SUBSYSTEM=="drm" RUN+="/bin/bash /opt/sys/vga_switch.sh"


